# systemd/claude-telegram-bot-healthcheck.service
[Unit]
Description=Healthcheck Claude Telegram Bot

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/claude-telegram-bot-healthcheck.sh
